# main/BizarroHomer.service
[Unit]
Description=Robot Control

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/BizarroHomerMain
Restart=always
RestartSec=0
KillMode=process
StandardOutput=file:/var/frc1511/BizarroHomer/dashboard/main_output.log
StandardError=file:/var/frc1511/BizarroHomer/dashboard/main_error.log

[Install]
WantedBy=multi-user.target
